# ansible/roles/aladdin-vault-app/files/etc/systemd/aladdin-initial-setup.timer
[Unit]
Description=Wait 5 minutes before starting aladdin-initial-setup service

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
