# ansible/roles/ipfs/files/etc/supervisord.conf

[program:ipfsoffline]
command=/usr/bin/ipfs daemon --offline
autostart=true
autorestart=true
priority=200
